# rc.local
#!/bin/sh
# This is run by SYS before other services are started
