# services/agetty@/run
#!/bin/sh
[ -r /etc/login.local ] && . /etc/login.local
[ "$1" = tty1 ] && options=--noclear
exec agetty $options -- "$1" 38400 linux
