# services/mdevd/run
#!/bin/sh
exec >/dev/null 2>&1
exec mdevd -C -D 3 -O 4
